# opt/nextcloudpi/build/pideploy/notify_push
#!/bin/sh
### BEGIN INIT INFO
# Provides:          notify_push
# Required-Start:    mariadb redis-server apache2
# Required-Stop:     
# Should-Start:      
# Should-Stop:       
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Push Notification Service
# Description:       Push Notification Service
### END INIT INFO

PATH=/sbin:/bin:/usr/sbin:/usr/bin
DESC="Push Notification"
NAME="notify_push"
DAEMON="/usr/local/etc/$NAME"
SCRIPTNAME="/etc/init.d/$NAME"

chmod +x /var/www/nextcloud/apps/notify_push/bin/aarch64/notify_push
. /lib/lsb/init-functions


case "$1" in
    start)
        log_daemon_msg "Starting $DESC" "$NAME"
	systemctl start notify_push
        log_end_msg $?
        ;;
    stop)
        log_daemon_msg "Stopping $DESC" "$NAME"
	pkill -f /var/www/nextcloud/apps/notify_push/bin/aarch64/notify_push
        log_end_msg $?
        ;;
 restart)
        killall notify_push &>/dev/null
        log_daemon_msg "Restarting $DESC" "$NAME"
	pkill -f /var/www/nextcloud/apps/notify_push/bin/aarch64/notify_push
        systemctl start notify_push
        log_end_msg $?
        ;;
    *)
        echo "Usage: $SCRIPTNAME {start|stop|restart}" >&2
        exit 3
        ;;
esac

exit 0

