# etc/init.d/crond
# vim /etc/rc.d/init.d/crond
#!/bin/bash
#
# description: crond daemon
# chkconfig: 2345 10 90
#
lockfile=/var/lock/subsys/crond
pidfile=/var/run/crond.pid
binary=/sbin/crond
[ -r /etc/rc.d/init.d/functions ] && . /etc/rc.d/init.d/functions
[ -r /etc/sysconfig/crond ] && . /etc/sysconfig/crond

start() {
    if [ -e $lockfile ]; then
        echo -n "crond daemon is already running: "
        success
        echo
    exit 0
fi
    echo -n "Starting crond:                "
    daemon --pidfile="$pidfile" $binary
    RETVAL=$?
    echo
if [ $RETVAL -eq 0 ]; then
    touch $lockfile
    return 0
else
    rm -f $lockfile $pidfile
    return 1
fi
}
stop() {
    if [ ! -e $lockfile ]; then
        echo -n "crond service is stopped: "
        success
        echo
        exit 1
    fi
        echo -n "Stopping crond daemon: "
        pkill -HUP crond
        RETVAL=$?
        echo
                                                                                                                                                                                                                                                                                                                                                                                                                                                                
    if [ $RETVAL -eq 0 ]; then
        rm -f $lockfile $pidfile
        return 0
    else
        return 1
fi
}
status() {
    if [ -e $lockfile ]; then
        echo "crond is running..."
    else
        echo "crond is stopped..."
    fi
}
usage() {
    echo "Usage: crond {start|stop|restart|status}"
}
case $1 in
start)
    start ;;
stop)
    stop ;;
restart)
    stop
    start
    ;;
status)
    status
    ;;
*)
    usage
    ;;
esac
