# etc/init.d/lighttpd
# vim /etc/rc.d/init.d/lighttpd
#!/bin/bash
#
# description: lighttpd daemon
# chkconfig: 2345 66 33
#
lockfile=/var/lock/subsys/lighttpd
pidfile=/var/run/lighttpd.pid
binary=/sbin/lighttpd
[ -r /etc/rc.d/init.d/functions ] && . /etc/rc.d/init.d/functions
[ -r /etc/sysconfig/lighttpd ] && . /etc/sysconfig/lighttpd

start() {
    if [ -e $lockfile ]; then
        echo -n "lighttpd daemon is already running: "
        success
        echo
    exit 0
fi
    echo -n "Starting lighttpd: "
    daemon --pidfile="$pidfile" $binary -f /etc/lighttpd/lighttpd.conf
    RETVAL=$?
    echo
if [ $RETVAL -eq 0 ]; then
    touch $lockfile
    return 0
else
    rm -f $lockfile $pidfile
    return 1
fi
}
stop() {
    if [ ! -e $lockfile ]; then
        echo -n "lighttpd service is stopped: "
        success
        echo
        exit 1
    fi
        echo -n "Stopping lighttpd daemon: "
        pkill -9 -u lighttpd
        RETVAL=$?
        echo
                                                                                                                                                                                                                                                                                                                                                                                                                                                                
    if [ $RETVAL -eq 0 ]; then
        rm -f $lockfile $pidfile
        return 0
    else
        return 1
fi
}
status() {
    if [ -e $lockfile ]; then
        echo "lighttpd is running..."
    else
        echo "lighttpd is stopped..."
    fi
}
usage() {
    echo "Usage: lighttpd {start|stop|restart|status}"
}
case $1 in
start)
    start ;;
stop)
    stop ;;
restart)
    stop
    start
    ;;
status)
    status
    ;;
*)
    usage
    ;;
esac

