# etc/init.d/xrdp
# vim /etc/rc.d/init.d/xrdp
#!/bin/bash
#
# description: xrdp daemon
# chkconfig: 2345 66 33
#
lockfile=/var/lock/subsys/xrdp
pidfile=/var/run/xrdp.pid
binary=/sbin/xrdp
cd /etc/xrdp
[ -r /etc/rc.d/init.d/functions ] && . /etc/rc.d/init.d/functions
[ -r /etc/sysconfig/xrdp ] && . /etc/sysconfig/xrdp

start() {
    [ -f /etc/xrdp/rsakeys.ini ] || bash -c "xrdp-keygen XRDP auto ; openssl req -x509 -newkey rsa:2048 -nodes -keyout key.pem -out cert.pem -days 3650 -subj '/C=GB/ST=London/L=London/O=Android/OU=Linux Deploy/CN=pi.hole' ; chmod 400 key.pem ; chmod 444 cert.pem" 
    if [ -e $lockfile ]; then
        echo -n "xrdp daemon is already running: "
        success
        echo
    exit 0
fi
    echo -n "Starting xrdp:                 "
    /sbin/xrdp-sesman
    daemon --pidfile="$pidfile" $binary
    RETVAL=$?
    echo
if [ $RETVAL -eq 0 ]; then
    touch $lockfile
    return 0
else
    rm -f $lockfile $pidfile
    return 1
fi
}
stop() {
    if [ ! -e $lockfile ]; then
        echo -n "xrdp service is stopped: "
        success
        echo
        exit 1
    fi
        echo -n "Stopping xrdp daemon: "
        pkill xrdp-sesman 
        pkill xrdp
        RETVAL=$?
        echo
                                                                                                                                                                                                                                                                                                                                                                                                                                                                
    if [ $RETVAL -eq 0 ]; then
        rm -f $lockfile $pidfile
        return 0
    else
        return 1
fi
}
status() {
    if [ -e $lockfile ]; then
        echo "xrdp is running..."
    else
        echo "xrdp is stopped..."
    fi
}
usage() {
    echo "Usage: xrdp {start|stop|restart|status}"
}
case $1 in
start)
    start ;;
stop)
    stop ;;
restart)
    stop
    start
    ;;
status)
    status
    ;;
*)
    usage
    ;;
esac
