# etc/update-motd.d/10-pihole
[ -f /usr/local/bin/pihole ] && echo && pihole status && echo && pihole version && echo
